# etc/supervisord.d/logstash.ini
[program:logstash]
command=/opt/qnib/logstash/bin/start.sh
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
